# deployment-examples/systemd/juicefs-llm-access-usage.service.template
[Unit]
Description=JuiceFS mount for StaticFlow llm-access usage analytics
After=network-online.target juicefs-llm-access.service
Wants=network-online.target
Requires=juicefs-llm-access.service

[Service]
Type=simple
User=ts_user
Group=ts_user
EnvironmentFile=/home/ts_user/.config/staticflow/llm-access-juicefs.env
ExecStartPre=/usr/bin/test -d /mnt/llm-access-usage
ExecStartPre=/usr/bin/mkdir -p /mnt/llm-access/llm-access-usage
ExecStart=/home/ts_user/.local/bin/juicefs mount --cache-dir /var/cache/juicefs/llm-access-usage --cache-size 40960 --log /var/log/juicefs/llm-access-usage.log --no-syslog --no-usage-report --writeback --subdir /llm-access-usage ${JUICEFS_META_URL} /mnt/llm-access-usage
ExecStop=/usr/bin/fusermount3 -uz /mnt/llm-access-usage
Restart=always
RestartSec=3
TimeoutStopSec=30
TimeoutStartSec=180
ExecStartPost=/usr/local/bin/staticflow-wait-llm-access-state mount-only
MemoryAccounting=yes
MemoryHigh=1800M
MemoryMax=2560M
MemorySwapMax=0
TasksMax=256
OOMPolicy=kill

[Install]
WantedBy=multi-user.target
